# systemd/tinyclaw.service
[Unit]
Description=TinyClaw Docker Compose Stack
After=docker.service
Requires=docker.service

[Service]
Type=simple
User=__USER__
WorkingDirectory=__WORKING_DIR__
ExecStart=/usr/bin/docker compose up
ExecStop=/usr/bin/docker compose down
Restart=always
RestartSec=10
TimeoutStartSec=120
TimeoutStopSec=60

StandardOutput=journal
StandardError=journal
SyslogIdentifier=tinyclaw

[Install]
WantedBy=multi-user.target
