# systemd/tinyclaw-queue.service
[Unit]
Description=TinyClaw Queue Processor
After=network.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__WORKING_DIR__
ExecStart=/usr/bin/node dist/queue-processor.js
Restart=always
RestartSec=5
Environment=NODE_ENV=production

StandardOutput=journal
StandardError=journal
SyslogIdentifier=tinyclaw-queue

[Install]
WantedBy=multi-user.target
